# server-config/tileserver-gl.service
[Unit]
Description=tileserver-gl
After=network.target

[Service]
After=tileserver-gl.service
ExecStart=tileserver-gl -c /mnt/maps/mapbh/tile-config.json
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
